# service/privoxy/run
#!/bin/bash

privoxy --no-daemon privoxy.cfg
